# overlay/etc/local.d/novaos.start
#!/bin/sh
# Automatic NovaOS Desktop startup script

# Start D-Bus and LightDM graphical display manager
rc-service dbus start 2>/dev/null || true
rc-service lightdm start 2>/dev/null || true
