# overlay/etc/local.d/novaos.start
#!/bin/sh

if ! id nova >/dev/null 2>&1; then
    adduser -D -g "NovaOS User" nova
    echo "nova:nova" | chpasswd
    addgroup nova wheel
    addgroup nova video
    addgroup nova input
fi

rc-service dbus start
rc-service elogind start
rc-service lightdm start
